# engine/src/main/unix/scripts/postinst
MODULE="streaming"

cat <<'EOF' > $MODULE
#!/bin/sh
### BEGIN INIT INFO
# Provides:          Stratio $MODULE
# Required-Start:    
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Description:       Stratio $MODULE startup script
### END INIT INFO
 
MODULE="streaming"
RUNAS=stratio

SCRIPT_ACTION=/opt/sds/$MODULE/bin/runActionEngine
SCRIPT_DATA=/opt/sds/$MODULE/bin/runDataEngine
SCRIPT_PROCESS=/opt/sds/$MODULE/bin/runProcessEngine
 
PIDFILE_ACTION=/var/run/sds/$MODULE_action.pid
LOGFILE_ACTION=/var/log/sds/$MODULE_action.log

PIDFILE_DATA=/var/run/sds/$MODULE_data.pid
LOGFILE_DATA=/var/log/sds/$MODULE_data.log

PIDFILE_PROCESS=/var/run/sds/$MODULE_process.pid
LOGFILE_PROCESS=/var/log/sds/$MODULE_process.log
 
start() {
  if [ -f $PIDFILE_ACTION ] && kill -0 $(cat $PIDFILE_ACTION); then
    echo 'Service streaming_action already running' >&2
  else
    echo 'Starting service streaming_action…' >&2
    local CMD="$SCRIPT_ACTION &> \"$LOGFILE_ACTION\" & echo \$!"
    su -c "$CMD" $RUNAS > "$PIDFILE_ACTION"
    echo 'Service streaming_action started' >&2
  fi
  if [ -f $PIDFILE_DATA ] && kill -0 $(cat $PIDFILE_DATA); then
    echo 'Service streaming_data already running' >&2
  else
    echo 'Starting service streaming_data…' >&2
    local CMD="$SCRIPT_DATA &> \"$LOGFILE_DATA\" & echo \$!"
    su -c "$CMD" $RUNAS > "$PIDFILE_DATA"
    echo 'Service streaming_data started' >&2
  fi
  if [ -f $PIDFILE_PROCESS ] && kill -0 $(cat $PIDFILE_PROCESS); then
    echo 'Service streaming_process already running' >&2
  else
    echo 'Starting service streaming_data…' >&2
    local CMD="$SCRIPT_PROCESS &> \"$LOGFILE_PROCESS\" & echo \$!"
    su -c "$CMD" $RUNAS > "$PIDFILE_PROCESS"
    echo 'Service streaming_process started' >&2
  fi
}
 
stop() {
  if [ ! -f "$PIDFILE_ACTION" ] || ! kill -0 $(cat "$PIDFILE_ACTION"); then
    echo 'Service streaming_action not running' >&2
  else
    echo 'Stopping service streaming_action…' >&2
    kill -15 $(cat "$PIDFILE_ACTION") && rm -f "$PIDFILE_ACTION"
    echo 'Service streaming_action stopped' >&2
  fi
  if [ ! -f "$PIDFILE_DATA" ] || ! kill -0 $(cat "$PIDFILE_DATA"); then
    echo 'Service streaming_data not running' >&2
  else
    echo 'Stopping service streaming_data…' >&2
    kill -15 $(cat "$PIDFILE_DATA") && rm -f "$PIDFILE_DATA"
    echo 'Service streaming_data stopped' >&2
  fi
  if [ ! -f "$PIDFILE_PROCESS" ] || ! kill -0 $(cat "$PIDFILE_PROCESS"); then
    echo 'Service streaming_process not running' >&2
  else
    echo 'Stopping service streaming_process…' >&2
    kill -15 $(cat "$PIDFILE_PROCESS") && rm -f "$PIDFILE_PROCESS"
    echo 'Service streaming_process stopped' >&2
  fi
}

status() {
	if [ ! -f "$PIDFILE_ACTION" ]; then
		echo 'Service streaming_action not running' >&2
	else
		echo 'Service streaming_action running with pid:' $(cat "$PIDFILE_ACTION") >&2
	fi
	if [ ! -f "$PIDFILE_DATA" ]; then
        echo 'Service streaming_data not running' >&2
    else
        echo 'Service streaming_data running with pid:' $(cat "$PIDFILE_DATA") >&2
    fi
    if [ ! -f "$PIDFILE_PROCESS" ]; then
        echo 'Service streaming_process not running' >&2
    else
        echo 'Service streaming_process running with pid:' $(cat "$PIDFILE_PROCESS") >&2
    fi
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  status)
    status
    ;;
  restart)
    stop
    start
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart}"
    exit 1
    ;;
esac

EOF

USER="stratio"
GROUP="stratio"

mv $MODULE /etc/init.d/$MODULE
chmod a+x /etc/init.d/$MODULE
mkdir -p /var/run/sds
mkdir -p /var/log/sds

if ! getent group $GROUP >/dev/null 2>&1; then
  groupadd $GROUP
fi

if ! getent passwd $USER >/dev/null 2>&1; then
  useradd -M -d /opt/sds -g $GROUP -s /bin/bash $USER
else
  usermod -G $GROUP $USER
fi

chown -R $USER:$GROUP /opt/sds/$MODULE
chmod a+x /opt/sds/$MODULE/bin/runActionEngine
chmod a+x /opt/sds/$MODULE/bin/runDataEngine
chmod a+x /opt/sds/$MODULE/bin/runProcessEngine
chown -R $USER:$GROUP /etc/sds/$MODULE
chown $USER:$GROUP /var/run/sds
chown $USER:$GROUP /var/log/sds

#/etc/init.d/$MODULE start
